# scripts/events-curator-deploy.timer
[Unit]
Description=Poll origin for new events-curator commits

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
