# scripts/events-curator-deploy.timer
[Unit]
Description=Poll origin for new events-curator commits

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
